# overlay/etc/cron.daily/update-tor-ipset
#!/bin/bash
/usr/local/bin/update-tor-ipset
service firewall reload-tor
